# trunk/etc/cron.monthly/delete_session_temp.cron
#!/bin/sh

rm -f /tmp/sess_* > /dev/null 2>&1

