# etc/cron.daily/delete_orphan_ivr.cron
#!/bin/sh

rm -f /var/lib/asterisk/sounds/*ivrrecording.wav > /dev/null 2>&1
