# trunk/etc/cron.daily/delete_cticalls_keys.cron
#!/bin/sh

if pgrep asterisk > /dev/null 2>&1
then

	/usr/sbin/asterisk -rx 'database deltree CTICALLS' > /dev/null 2> /dev/null

fi

