# trunk/etc/cron.daily/delete_fax.cron
#!/bin/sh

find /var/spool/asterisk/fax -type f -mtime +15 -name "*.*" -exec rm -f {} \;
