# trunk/etc/cron.daily/delete_recording.cron
#!/bin/sh

find /var/spool/asterisk/monitor -type f -mtime +15 -name "*.*" -exec rm -f {} \;
